# agents/sample-agent/workspace/skills/podcast-summary/scripts/launchagents/com.ironclaw.whisper-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ironclaw.whisper-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>__WHISPER_DIR__/build/bin/whisper-server</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>__PORT__</string>
        <string>--model</string>
        <string>__WHISPER_DIR__/models/ggml-small.en.bin</string>
        <string>--convert</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/whisper-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/whisper-bridge.log</string>

    <key>WorkingDirectory</key>
    <string>__WHISPER_DIR__</string>
</dict>
</plist>
