# scripts/launchagents/com.ironclaw.health-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Install: cp scripts/launchagents/com.ironclaw.health-watch.plist ~/Library/LaunchAgents/
     Load:    launchctl load ~/Library/LaunchAgents/com.ironclaw.health-watch.plist
     Unload:  launchctl unload ~/Library/LaunchAgents/com.ironclaw.health-watch.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ironclaw.health-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeff/ironclaw/scripts/watch-health-import.sh</string>
    </array>

    <!-- Poll every 5 minutes. More reliable than WatchPaths on iCloud Drive
         (a network-backed filesystem where WatchPaths is documented as flaky). -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jeff/Library/Logs/ironclaw/health-watch.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jeff/Library/Logs/ironclaw/health-watch.log</string>
</dict>
</plist>
